# deploy/systemd/chat-grid.service
[Unit]
Description=Endiginous signaling server
After=network.target

[Service]
Type=simple
# This file is a reference example.
# Use deploy/scripts/install_service.sh to generate an environment-specific unit.
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
